# .bash/profile.d/convenience
#!/bin/bash

function reload {
  rvm reload
  source "$HOME/.bash_profile"
}
