# .bash/profile.d/convenience
#!/bin/bash

function reload {
  source "$HOME/.bash_profile"
  rvm reload
}
