# .bash/profile.d/path
#!/bin/bash

PATH="$HOME/.bash/bin:$PATH:$HOME/bin
