# .bash/profile.d/path
#!/bin/bash

PATH="$HOME/.bash/bin:$HOME/.bash/bin/ssl:$PATH:$HOME/work/git-utils/bin"
