# .bash/profile.d/x11
function x11 {
  ( bash startx & ) &
}

