# .bash/profile.d/x11
function x11 {
  ( setsid bash startx & ) &
}

